# recipes/emss-rsc/emss-rsc/rsc.service
[Unit]
Description=EMSS Remote System Controller (RSC) service
After=network.target
RequiresMountsFor=/media/card

[Service]
WorkingDirectory=/home/rsc
ExecStart=/usr/sbin/rsc
Restart=always
RestartSec=2500ms
;User=rsc

[Install]
WantedBy=multi-user.target
